# deb/DEBIAN/postinst
#!/bin/sh -e
#
#
#

systemctl daemon-reload
systemctl enable ipchangecheck.service
/etc/cron.daily/motdUpdate > /dev/null 2>&1
