# deb/DEBIAN/postinst
#!/bin/sh -e
#
#
#

sysctl --system
systemctl daemon-reload
systemctl disable ipchangecheck.service
systemctl enable ipchangecheck.service
/etc/cron.daily/motdUpdate > /dev/null 2>&1
