# deb/DEBIAN/prerm
#!/bin/sh -e
systemctl disable ipchangecheck.service
