# deb/DEBIAN/prerm
#!/bin/sh
systemctl disable ipchangecheck.service