# deb/etc/systemd/system/docker-proxy-environment.service
[Service]
ExecStart=
ExecStart=/etc/systemd/system/docker-proxy-environment.service.d/start-script.sh
Restart=no
